# .setup/spelltracker.conf
[program:spelltracker]
command=gunicorn spelltracker.wsgi
autostart=yes
directory=/home/adam/spelltracker
startretries=3
user=adam
stderr_logfile=/var/log/gunicorn/gunicorn.err.log
stdout_logfile=/var/log/gunicorn/gunicorn.out.log
